# profile.d/brew
#!/bin/bash

source ~/.bash/functions

if [ -x brew ] && [ -f $(brew --prefix)/etc/bash_completion ]; then
  . $(brew --prefix)/etc/bash_completion
fi
