# profile.d/docker-machine
#!/bin/bash

source ~/.bash/functions

eval $(docker-machine env default 2>/dev/null)
